# content/service/xray/run
#!/bin/sh

export XRAY_LOCATION_ASSET=/config

# Run Xray
if [[ "${CONF_FILE_NAME}" = "" ]]; then
    exec xray run -confdir /config
else
    exec xray run -c /config/"${CONF_FILE_NAME}"
fi 